# ops/polymarket-weather-bot/systemd/polymarket-weather-monitor.service
[Unit]
Description=Polymarket Weather Paper Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kai/.openclaw/workspace/ops/polymarket-weather-bot

# Risk profile (can be overridden with systemctl --user edit ...)
Environment=INTERVAL_SEC=300
Environment=MIN_HOURS_TO_EXPIRY=0
Environment=MAX_POSITIONS_PER_CITY=2
Environment=MAX_EVENT_CLUSTER_EXPOSURE_USD=100
Environment=EXIT_EDGE_FLOOR=0.01
Environment=MIN_HOLDING_MINUTES_FOR_EDGE_EXIT=10
Environment=CONFIRM_TICKS=2
Environment=TRADE_SIZE_USD=100
Environment=MAX_OPEN_EXPOSURE_USD=200
Environment=DAILY_STOP_LOSS_USD=-50
Environment=DAILY_NEW_OPEN_NOTIONAL_CAP_USD=0
Environment=PAPER_BANKROLL_USD=1000
Environment=KELLY_FRACTION_CORE=0.20
Environment=KELLY_FRACTION_TAIL=0.08
Environment=MAX_BET_FRACTION=0.10
Environment=TAIL_SIZE_CAP_FRACTION=0.5
Environment=MIN_EDGE_FOR_ENTRY=0.02
Environment=MIN_OPEN_SIZE_USD=10
Environment=ROBUSTNESS_MU_SHIFT_C=0.7
Environment=ROBUSTNESS_SIGMA_SCALE_LOW=0.85
Environment=ROBUSTNESS_SIGMA_SCALE_HIGH=1.15
Environment=ROBUSTNESS_MIN_EDGE=0.0
Environment=ENABLE_EDGE_ROTATION=1
Environment=ROTATION_MIN_EDGE_DELTA=0.05
Environment=ROTATION_MIN_EV_PER_USD_DELTA=0.08
Environment=ROTATION_MIN_HOLDING_MINUTES=10
Environment=MAX_ROTATIONS_PER_RUN=1
Environment=ROTATION_REQUIRE_PROFIT=1
Environment=COMPOUND_ENABLED=1
Environment=COMPOUND_TRADE_SIZE_FRACTION=0.10
Environment=COMPOUND_MAX_OPEN_EXPOSURE_FRACTION=0.20
Environment=COMPOUND_DAILY_STOP_LOSS_FRACTION=0.03
Environment=COMPOUND_TRADE_SIZE_MIN_USD=10
Environment=COMPOUND_TRADE_SIZE_MAX_USD=100
Environment=COMPOUND_MAX_OPEN_EXPOSURE_MIN_USD=200
Environment=COMPOUND_MAX_OPEN_EXPOSURE_MAX_USD=300
Environment=COMPOUND_DAILY_STOP_LOSS_MIN_ABS_USD=50
Environment=COMPOUND_DAILY_STOP_LOSS_MAX_ABS_USD=50

ExecStart=/home/kai/.openclaw/workspace/ops/polymarket-weather-bot/scripts/monitor_loop.sh

Restart=always
RestartSec=10

[Install]
WantedBy=default.target
